# archlinux/i3/wallpaper.cron
* * * * * HOME/.config/i3/desktop.sh; sleep 20; HOME/.config/i3/desktop.sh; sleep20; HOME/.config/i3/desktop.sh; sleep 20; HOME/.config/i3/desktop.sh;